# deploy/sewing-web-backup.timer
[Unit]
Description=Daily Shagaem Vmeste database backup

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=10m
Unit=sewing-web-backup.service

[Install]
WantedBy=timers.target
